# scripts/tinyerp_demo/crontab
# m h dom mon dow user  command
0 0 * * * root /opt/tinyerp/create_demo_db
